# debian-live/config/includes.chroot/etc/systemd/system/determine_mode.service
[Unit]
Description=determine non-persistent or persistent mode

[Service]
ExecStart=/opt/iksdp/bin/determine_mode.sh

[Install]
WantedBy=multi-user.target