# debian-live/config/includes.chroot/etc/systemd/system/determinemode.service
[Unit]
Description=determine non-persistent or persistent mode

[Service]
ExecStart=/opt/iksdp/bin/determinemode.sh

[Install]
WantedBy=multi-user.target